# usr/lib/systemd/system/pve-agent.service
[Unit]
Description=Monitor Proxmox VE virtual machines

[Service]
Type=simple
ExecStart=/usr/bin/pve-agent